# packages/dpkg/vdu_1.0-1/etc/systemd/user/vdu.service
[Unit]
Description=VDU App Fuse Service

[Service]
Type=forking
ExecStartPre=/bin/sleep 30
ExecStart=/usr/local/bin/vdu-app-fuse /mnt/vdu
Restart=on-failure

[Install]
WantedBy=default.target
