# packages/vdu.service
[Unit]
Description=VDU App Fuse Service

[Service]
Type=forking
ExecStart=/home/user/vdu-app-fuse /mnt/vdu
Restart=on-failure

[Install]
WantedBy=multi-user.target
